# server/docker/supervisor/warriors-server.conf
[program:warriors-server]
command=nodejs /srv/warriors-server.js
user=nobody
restart=true
